# bootstrap.conf
[program:bootstrap]
command=gunicorn wsgi --bind 0.0.0.0:8000 --pid /tmp/bootstrap.pid
directory=/var/www/bootstrap
autostart=true
autorestart=true